# deploy/systemd/autotrader-pro.service
# /etc/systemd/system/autotrader-pro.service
[Unit]
Description=AUTOTRADER-PRO (8095)
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=roteiro_ds
WorkingDirectory=/home/roteiro_ds/AUTOTRADER-PRO

Environment=PORT=8095
Environment=DATA_DIR=/home/roteiro_ds/AUTOTRADER-PRO/data
Environment=DIST_PATH=/home/roteiro_ds/AUTOTRADER-PRO/dist

# valida o arquivo antes de subir (evita cair por erro de sintaxe)
ExecStartPre=/usr/bin/node --check /home/roteiro_ds/AUTOTRADER-PRO/server.js

ExecStart=/usr/bin/node /home/roteiro_ds/AUTOTRADER-PRO/server.js

Restart=always
RestartSec=2
TimeoutStopSec=10
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
