# deploy/systemd/autotrader-pro-worker.service
# /etc/systemd/system/autotrader-pro-worker.service
[Unit]
Description=AUTOTRADER-PRO worker (gera data/pro.json)
After=network.target

[Service]
Type=oneshot
User=roteiro_ds
WorkingDirectory=/home/roteiro_ds/AUTOTRADER-PRO
ExecStart=/usr/bin/python3 /home/roteiro_ds/AUTOTRADER-PRO/worker_pro.py
