# deploy/systemd/autotrader-pro-worker.service
[Unit]
Description=AUTOTRADER-PRO Worker (Python) - gera JSON
After=network.target

[Service]
Type=oneshot
User=roteiro_ds
WorkingDirectory=/home/roteiro_ds/AUTOTRADER-PRO
Environment=DATA_DIR=/home/roteiro_ds/AUTOTRADER-PRO/data
ExecStart=/usr/bin/python3 -m engine.worker_pro
